# tools/sovereign-hive@.service
[Unit]
Description=Sovereign Hive Strategy: %i
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=ec2-user
Group=ec2-user
WorkingDirectory=/app/sovereign-hive
Environment="PATH=/app/sovereign-hive/venv/bin:/usr/local/bin:/usr/bin"
Environment="PYTHONUNBUFFERED=1"
Environment="STRATEGY_FILTER=%i"

# Load secrets from Secrets Manager (in-memory only, not on disk)
ExecStartPre=/app/sovereign-hive/load_secrets.sh
EnvironmentFile=/run/sovereign-hive/env

ExecStart=/app/sovereign-hive/venv/bin/python3 sovereign_hive/run_simulation.py --live
Restart=always
RestartSec=30

StandardOutput=append:/var/log/sovereign-hive/%i.log
StandardError=append:/var/log/sovereign-hive/%i_error.log

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/var/log/sovereign-hive /app/sovereign-hive /tmp /run/sovereign-hive

[Install]
WantedBy=multi-user.target
