# tools/sovereign-hive@.service
[Unit]
Description=Sovereign Hive Strategy: %i
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=ec2-user
Group=ec2-user
WorkingDirectory=/app/hft
Environment="PATH=/app/hft/venv/bin:/usr/local/bin:/usr/bin"
Environment="PYTHONUNBUFFERED=1"
Environment="STRATEGY_FILTER=%i"
EnvironmentFile=/app/hft/.env
ExecStart=/app/hft/venv/bin/python3 sovereign_hive/run_simulation.py --live
Restart=always
RestartSec=30
StandardOutput=append:/var/log/quesquant/sovereign_hive_%i.log
StandardError=append:/var/log/quesquant/sovereign_hive_%i_error.log

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/var/log/quesquant /app/hft /tmp

[Install]
WantedBy=multi-user.target
